# templates/service/ccs.service.epp
<%- | String $desc,
      String $user,
      String $group,
      String $cmd,
      String $workdir
| -%>
[Unit]
Description=<%= $desc %>
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=<%= $workdir %>
ExecStart=<%= $cmd %>
Restart=on-failure
RestartSec=42s
User=<%= $user %>
Group=<%= $group %>

[Install]
WantedBy=multi-user.target
