# templates/service/status-email-user@.service.epp
<%- | String $envfile | -%>
[Unit]
Description=status email for %i

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/systemd-email $EMAIL %i
EnvironmentFile=<%= $envfile %>
